# systemd/patch-pam-alternative-passwords.path
[Unit]
Description=Watch PAM files for alternative password configuration updates

[Path]
PathChanged=/etc/pam.d/system-auth
PathChanged=/etc/pam.d/su
PathChanged=/etc/pam.d/su-l
PathChanged=/etc/pam.d/vlock
PathChanged=/etc/xdg/quickshell/caelestia/assets/pam.d/passwd
Unit=patch-pam-alternative-passwords.service

[Install]
WantedBy=multi-user.target
